# cronjobs/cron.daily/syncronize_package_databases.sh
#! /bin/bash
# This script will help maintain the package databases by sycronizing them. 

echo "Syncronizing package databases..."
sudo pacman -Sy 
